# agents/com.workspacemanager.janitor.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Janitor (feature 1). Weekly DRY-RUN scan — writes a review report, moves
  nothing. Review it, then run `workspace-manager janitor --apply` by hand if
  you want the flagged items relocated. Runs Sundays at 09:00 local.
  install.sh substitutes __WM_BIN__ and __LOG_DIR__.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.workspacemanager.janitor</string>
    <key>ProgramArguments</key>
    <array>
        <string>__WM_BIN__</string>
        <string>janitor</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>__LOG_DIR__/janitor.log</string>
    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/janitor.log</string>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
