# agents/com.workspacemanager.downloadsorter.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Download sorter (feature 2). Fires whenever ~/Downloads changes, then
  debounces via ThrottleInterval so a burst of files triggers one run.
  install.sh substitutes __WM_BIN__, __DOWNLOADS__, and __LOG_DIR__.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.workspacemanager.downloadsorter</string>
    <key>ProgramArguments</key>
    <array>
        <string>__WM_BIN__</string>
        <string>sort</string>
    </array>
    <key>WatchPaths</key>
    <array>
        <string>__DOWNLOADS__</string>
    </array>
    <key>ThrottleInterval</key>
    <integer>30</integer>
    <key>StandardOutPath</key>
    <string>__LOG_DIR__/downloadsorter.log</string>
    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/downloadsorter.log</string>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
